# configs/systemd/cmux-x11vnc.service
[Unit]
Description=Cmux x11vnc bridge
After=cmux-xvfb.service
Requires=cmux-xvfb.service

[Service]
Type=simple
Environment=DISPLAY=:1
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/usr/bin/x11vnc -display :1 -forever -shared -rfbport 5901 -localhost -nopw -quiet
Restart=always
RestartSec=3
StandardOutput=append:/var/log/cmux/x11vnc.log
StandardError=append:/var/log/cmux/x11vnc.log

[Install]
WantedBy=cmux.target
